# deploy/cgbrains.service
[Unit]
Description=cgbrains – AI-powered query builder
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/cgbrains
EnvironmentFile=/var/www/cgbrains/.env
ExecStart=/var/www/cgbrains/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8100 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
